# deploy/ringdown-mcp.service
[Unit]
Description=ringdown-mcp — Entra-gated control plane (query + rule/target CRUD)
After=network-online.target
Wants=network-online.target

[Service]
User=ringdown
Group=ringdown
WorkingDirectory=/opt/ringdown
# Binds 127.0.0.1:8787 (nginx fronts it with the public host + TLS). No
# privileged port, so NO capabilities are granted at all.
ExecStart=/opt/ringdown/.venv/bin/python -m ringdown.mcp_server
Restart=always
RestartSec=3
# Fast, clean restarts: cap the graceful-drain wait on long-lived streamable-HTTP
# clients at 5s (vs the 90s default), then SIGKILL; treat the signal exit as clean.
TimeoutStopSec=5
SuccessExitStatus=143 137

# --- sandbox (no capabilities needed) ---
NoNewPrivileges=yes
UMask=0077
CapabilityBoundingSet=
ProtectSystem=strict
ProtectHome=yes
PrivateTmp=yes
PrivateDevices=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX
RestrictNamespaces=yes
LockPersonality=yes
MemoryDenyWriteExecute=yes
ReadWritePaths=/opt/ringdown/var

[Install]
WantedBy=multi-user.target
